# cron/crontab
MAILTO=webmaster@lojban.org
37 */5 * * * /home/sphaproxy/haproxy/rebuild_images.sh
1  1   * * 1 /home/sphaproxy/haproxy/renew_certs.sh
