# cron/crontab
MAILTO=webmaster@lojban.org
37 */5 * * * /home/sphaproxy/haproxy/rebuild_images.sh
