# cron/crontab
MAILTO=webmaster@lojban.org
37 * * * * /home/sphaproxy/haproxy/rebuild_images.sh
